# systemd-work/my-one-short-task.service
[Unit]
    Description=Cleanup old Foo data

[Service]
    Type=oneshot
    ExecStart=/home/heirat/systemd-work/one-shot-service.bash