# systemd-work/my-one-short-task.timer
[Unit]
Description=Run every 60 sec

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w