# bash-scripts/backup-task.timer
[Unit]
    Description=Run every wednesday

[Timer]
    OnCalendar=Wed *-*-* 0:0:30
    Persistent=true